# systemd/healthcheck.timer
[Unit]
Description=Run ADE healthcheck daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
